# ddns.service.in
[Unit]
Description=ZX2C4 DDNS Server

[Service]
Type=simple
ExecStart=BINDIR/ddns serve
EnvironmentFile=SYSCONFDIR/ddns.conf
Restart=always
StandardInput=null
StandardOutput=journal
StandardError=journal
ProcSubset=pid
DynamicUser=true
NoNewPrivileges=true
LimitNOFILE=65535
ProtectSystem=strict
ProtectHome=true
StateDirectory=ddns
CacheDirectory=ddns
StateDirectoryMode=0700
CacheDirectoryMode=0700
PrivateTmp=true
PrivateDevices=true
ProtectHostname=true
ProtectClock=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6
LockPersonality=true
RestrictNamespaces=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
RestrictSUIDSGID=true
RemoveIPC=true
SystemCallFilter=@basic-io @file-system @io-event @ipc @network-io @process @signal @sync @timer mprotect madvise mremap brk copy_file_range fadvise64 fadvise64_64 flock getcpu getrandom readdir sched_yield sched_getaffinity sendfile sendfile64 splice tee umask sysinfo uname ioctl
